# workspace/crontab/laradock
* * * * * laradock cd /var/www/mid9songposts && /usr/bin/php /var/www/mid9songposts/artisan schedule:run >> /dev/null 2>&1
